# platform_restart/pr-custom.target
[Unit]
Description=Custom action script for platform restart
DefaultDependencies=no
After=cleanup-tmp.service
Before=do-return.service
